# postfix/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -f /etc/rsyslog.conf -n

[program:postfix]
command=/bin/postfix
